# templates/operator/systemd-user/edc-weekly-snapshot.service
[Unit]
Description=edc weekly snapshot

[Service]
Type=oneshot
EnvironmentFile=%h/.config/edcloud/edc.env
ExecStart=%h/.local/bin/edc snapshot -d weekly-snapshot
